# ops/backup/surau-alert@.service
[Unit]
Description=Telegram alert for failed unit %i

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/surau-alert %i
