# ops/backup/surau-backup-watchdog.timer
[Unit]
Description=Run Surau backup dead-man check every 2 hours

[Timer]
OnCalendar=00/2:15:00
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
